# centos-8/overlay/etc/systemd/system/console-getty.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/usr/sbin/agetty -o '-p -- \\u' --noclear --keep-baud --autologin root console 115200,38400,9600 $TERM